# systemd/backup-mysql.service
[Unit]
Description=Backup MySQL POS

[Service]
Type=oneshot
ExecStart=/opt/pos/bin/backup-mysql.sh
